# phpdocker/php-fpm/supervisord.conf
[supervisord]
user=root
nodaemon=true
logfile=/application/storage/logs/supervisord.log

[unix_http_server]
file=/var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
user=root
serverurl=unix:///var/run/supervisor.sock

[program:php-fpm]
user=root
command=/usr/sbin/php-fpm7.4 -O
autostart=true
autorestart=true

[program:php-websocket]
user=root
command=php /application/artisan websockets:serve
autostart=true
autorestart=true
stdout_logfile=/application/storage/logs/websocket.log